# templates/etc/supervisor.conf.in
[inet_http_server]
port = ${config:supervisord-port}
username = 
password = 

[supervisord]
logfile = ${buildout:directory}/var/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = ${buildout:directory}/var/supervisord.pid
nodaemon = false

[supervisorctl]
serverurl = http://localhost:${config:supervisord-port}

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface


[program:memcached]
command = ${memcached:location}/bin/memcached -v -p ${config:memcached-port} -U ${config:memcached-port} -m ${config:memcached-memory}
process_name = memcached
directory = ${memcached:location}
priority = 10
redirect_stderr = false

[program:solr]
#command = java [-Djava.security.egd=file:/dev/urandom -jar start.jar] ${buildout:parts-directory}/solr true
command = ${buildout:bin-directory}/solr-instance fg
process_name = solr
directory = ${buildout:parts-directory}/solr
priority = 20
redirect_stderr = false

